# config/s6/start/finish
#!/bin/sh
# 退出码 125 告诉 s6-svscan 不要重启此服务（一次性初始化服务）
exit 125
